# voice-room.service
[Unit]
Description=Voice Room Application
After=network.target

[Service]
Type=simple
User=voice-room
WorkingDirectory=/opt/voice-room
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="HOST=127.0.0.1"
Environment="CORS_ORIGIN=https://aiternitas.ru"
EnvironmentFile=-/opt/voice-room/.env
ExecStart=/usr/bin/node server/server.mjs
ExecStartPre=/bin/sleep 3
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=voice-room

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/voice-room

[Install]
WantedBy=multi-user.target

